# deep_talk.service
[Unit]
Description=Deep Talk Telegram Bot
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/deep_talk
ExecStart=/home/deploy/deep_talk/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
